# manifest/coldfusion.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!--
	SMF manifest to manage a Coldfusion 7 server. (May work with other
	versions, I don't know.) Runs as root.
	
	start - start server
	stop - stop server
	restart - restart server
	
	R Fisher / SNLTD 05/2007
-->

<service_bundle type='manifest' name='export'>

	<service name='snltd/coldfusion' type='service' version='0'>

		<create_default_instance enabled='true'/>

		<single_instance/>

		<dependency name='filesystem'
			grouping='require_all'
			restart_on='none'
			type='service'>

			<service_fmri value='svc:/system/filesystem/local'/>
		</dependency>

		<dependency name='network'
			grouping='require_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/network/initial'/>
		</dependency>

		<exec_method name='start'
			type='method'
			exec='/usr/local/coldfusion/bin/coldfusion start'
			timeout_seconds='30'>
			<method_context>
				<method_credential user='root' group='root'/>
			</method_context>
		</exec_method>

		<exec_method name='stop'
			type='method'
			exec='/usr/local/coldfusion/bin/coldfusion stop'
			timeout_seconds='60'>
			<method_context>
				<method_credential user='root' group='root'/>
			</method_context>
		</exec_method>

		<exec_method name='restart'
			type='method'
			exec='/usr/local/coldfusion/bin/coldfusion restart'
			 timeout_seconds='60'>
			<method_context>
				<method_credential user='root' group='root'/>
			</method_context>
		</exec_method>

		<stability value='Unstable'/>

		<template>
			<common_name>
				<loctext xml:lang='C'>Coldfusion</loctext>
			</common_name>
		</template>

	</service>

</service_bundle>
