# manifest/mediatomb.xml
<?xml version='1.0'?>

<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!--
	SMF manifest for Mediatomb 
-->

<service_bundle type='manifest' name='export'>

	<service name='snltd/mediatomb' type='service' version='0'>

		<create_default_instance enabled='true'/>

		<single_instance/>

		<dependency
			name='filesystem'
			grouping='require_all'
			restart_on='none'
			type='service'>

			<service_fmri value='svc:/system/filesystem/local'/>

		</dependency>

		<dependency
			name='network'
			grouping='require_all'
			restart_on='none'
			type='service'>
				
			<service_fmri value='svc:/network/initial'/>

		</dependency>

		<exec_method
			name='start'
			type='method'
			exec='/usr/local/mediatomb/bin/mediatomb
				-d 
				-c /config/mediatomb/config.xml 
				-u mediatomb 
				-l /var/log/mediatomb/mediatomb.log' 
				timeout_seconds='30'>

			<method_context>
				<method_environment>
					<envvar name='LD_PRELOAD' value='/usr/lib/0@0.so.1'/>
				</method_environment>
			</method_context>

		</exec_method>

		<exec_method
			name='stop'
			type='method'
			exec=':kill'
			timeout_seconds='30'>

			<method_context>
				<method_credential user='mediatomb' group='mediatomb'/>
			</method_context>

		</exec_method>

		<stability value='Unstable'/>

		<template>
			<common_name>
				<loctext xml:lang='C'>MediaTomb</loctext>
			</common_name>
		</template>

	</service>

</service_bundle>
